# files/cart.service
[Unit]
Description = Cart Service
[Service]
User=roboshop
Environment=REDIS_HOST=dev-redis.kiranprav.link
Environment=CATALOGUE_HOST=dev-catalogue.kiranprav.link
ExecStart=/bin/node /app/server.js
SyslogIdentifier=cart

[Install]
WantedBy=multi-user.target